# debian/remotertrx-systemd.service
;;;;; Author: Richard Neese<kb3vgw@gmail.com>

[Unit]
Description=remotetrx svxlink remote connection server
After=network.target remote-fs.target

[Service]
; service
PIDFile=/run/svxlink/remotetrx.pid
RuntimeDirectory=svxlink
ExecStart=/usr/bin/remotetrx
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=on-failure
; exec
WorkingDirectory=/run/svxlink
User=svxlink
Group=svxlink

[Install]
WantedBy=multi-user.target
