# src/svxlink/systemd/remotetrx.service.in
;;;;; Author: Richard Neese<kb3vgw@gmail.com>

[Unit]
Description=svxlink repeater reflector software
Documentation="man:svxlink"
After=network.target remote-fs.target time.target sound.target

[Service]
Type=simple
User=svxlink
Group=daemon
EnvironmentFile=-/etc/default/svxreflector
Environment=CFGFILE=/etc/svxlink/svxreflector.conf
RuntimeDirectory=svxlink
ExecStart=/usr/bin/svxreflector --config $CFGFILE
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
