# src/svxlink/systemd/remotetrx.service.in
;;;;; Author: Richard Neese<kb3vgw@gmail.com>

[Unit]
Description=open-link repeater reflector software
Documentation="man:open-link"
After=network.target remote-fs.target time.target sound.target

[Service]
Type=simple
User=open-link
Group=daemon
EnvironmentFile=-/etc/default/svxreflector
Environment=CFGFILE=/etc/open-link/svxreflector.conf
RuntimeDirectory=open-link
ExecStart=/usr/bin/svxreflector --config $CFGFILE
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
