# src/svxlink/systemd/svxlink.service.in
;;;;; Author: Richard Neese<kb3vgw@gmail.com>

[Unit]
Description=open-link repeater control software
Documentation="man:open-link"
Requires=open-link_gpio_setup.service
After=network.target remote-fs.target time.target sound.target open-link_gpio_setup.service

[Service]
Type=simple
User=open-link
Group=daemon
EnvironmentFile=-/etc/default/open-link
Environment=CFGFILE=/etc/open-link/open-link.conf
RuntimeDirectory=open-link
ExecStart=/usr/bin/open-link --config $CFGFILE
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
