# src/svxlink/systemd/svxlink_gpio_setup.service.in
[Unit]
Description=SvxLink repeater control software GPIO setup
Documentation=man:svxlink(1)

[Service]
Type=oneshot
ExecStart=@SBIN_INSTALL_DIR@/svxlink_gpio_up
ExecStop=@SBIN_INSTALL_DIR@/svxlink_gpio_down
#Restart=on-failure
TimeoutStartSec=60
WorkingDirectory=@SVX_SYSCONF_INSTALL_DIR@
RemainAfterExit=True
