# src/svxlink/systemd/svxlink.service.in
;;;;; Author: Richard Neese<kb3vgw@gmail.com>

[Unit]
Description=svxlink repeater control software
Documentation="man:svxlink(1)"
Requires=svxlink_gpio_setup.service
After=network.target remote-fs.target time.target sound.target svxlink_gpio_setup.service

[Service]
Type=simple
User=svxlink
Group=daemon
EnvironmentFile=-/etc/default/svxlink
Environment=CFGFILE=/etc/svxlink/svxlink.conf
RuntimeDirectory=svxlink
ExecStart=/usr/bin/svxlink --config $CFGFILE --logfile /var/log/svxlink --pidfile /run/svxlink/svxlink.pid
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
