# src/svxlink/systemd/svxreflector.service.in
[Unit]
Description=SvxReflector SvxLink node linking software
Documentation=man:svxreflector(1)
After=network.target remote-fs.target syslog.target time.target

[Service]
EnvironmentFile=/etc/default/svxreflector
PIDFile=/run/svxreflector.pid
ExecStartPre=-/bin/touch /var/log/svxreflector
ExecStartPre=-/bin/chown $RUNASUSER /var/log/svxreflector
ExecStart=/bin/sh -c '@BIN_INSTALL_DIR@/svxreflector --logfile=/var/log/svxreflector --config=$CFGFILE --pidfile=/run/svxreflector.pid --runasuser=$RUNASUSER'
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
TimeoutStartSec=60
LimitCORE=infinity
WorkingDirectory=@SVX_SYSCONF_INSTALL_DIR@

[Install]
WantedBy=multi-user.target
